# robcutmore-com/config/gunicorn_upstart.conf
description "Gunicorn server for robcutmore.com"

start on net-device-up
stop on shutdown

respawn

setuid rc

exec /home/rc/projects/robcutmore/robcutmore-com/bin/start_server.sh